# hardenedbsd/secadm/pkg-plist
etc/rc.d/secadm
include/secadm.h
lib/libsecadm.so
lib/libsecadm.so.1
share/man/man5/secadm.rules.5.gz
share/man/man8/secadm.8.gz
sbin/secadm
